# rg353v-custom/board/rg353v/rootfs-overlay/etc/inittab
::sysinit:/sbin/game_init
::restart:/sbin/game_init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
